# bootstrap/playbooks/files/supervisor.conf
[program:{{name}}]
command={{cmd}}
redirect_stderr=true
stdout_logfile=/var/log/{{name}}.log
