# utils/solar-worker.conf
# solar-worker
#

description     "solar-worker upstart script"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [06]

# Prepare the environment
# Change ownership to the user running the process
pre-start script
    . /etc/default/solar-worker
    mkdir -p /var/log/solar
    mkdir -p /var/opt/solar
    chown -R $SOLAR_UID:$SOLAR_GID /var/log/solar
    chown -R $SOLAR_UID:$SOLAR_GID /var/opt/solar
end script

# If the process quits unexpectadly trigger a respawn
respawn

# Start the process
script
    . /etc/default/solar-worker
    export SOLAR_CONFIG_OVERRIDE=/.solar_config_override
    exec start-stop-daemon --start --chuid $SOLAR_UID --make-pidfile --pidfile $SOLAR_PIDFILE --exec /usr/bin/env solar-worker
end script
